# supervisord.conf
[program:celeryd]
command = python current/manage.py celeryd -l info -E --settings=settings_production
