# supervisord.conf
[program:celeryd]
command = python manage.py celeryd -l info -E --settings=settings_production
